# systemd/surgeprotector.timer
[Unit]
Description=Block Tor Exit traffic to flooded IP addresses

[Timer]
OnCalendar=*-*-* *:*

[Install]
WantedBy=timers.target
